# etc/systemd/user/midicontrol@.service
[Unit]
Description=Control midi keyboards
StopWhenUnneeded=true

[Service]
ExecStart=/usr/bin/midicontrol start -d %I
